# scripts/com.openrouter.log-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openrouter.log-export</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/pnpm</string>
        <string>run</string>
        <string>export-logs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/admin/Dev/Projects/openrouter-tracker</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/admin/Library/Logs/openrouter-export.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/Library/Logs/openrouter-export.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
